# roles/snapshots/files/grub-btrfs/systemd/grub-btrfs.path
# /usr/lib/systemd/system/grub-btrfs.path

[Unit]
Description=Monitors for new snapshots
DefaultDependencies=no
Requires=-.mount
After=-.mount
BindsTo=-.mount

[Path]
PathModified=/.snapshots

[Install]
WantedBy=-.mount
